# services/watchdog/tests/fixtures/systemd_unit_expected.service
[Unit]
Description=Gimle Palace agent watchdog (GIM-63)
After=network.target

[Service]
Type=simple
ExecStart=/path/to/.venv/bin/python -m watchdog run --config /home/user/.paperclip/watchdog-config.yaml
Restart=on-failure
RestartSec=10s
StandardOutput=append:/home/user/.paperclip/watchdog.log
StandardError=append:/home/user/.paperclip/watchdog.err

[Install]
WantedBy=default.target
